# systemd/clawpi.target
[Unit]
Description=ClawPi normal operating target
AllowIsolate=yes
StopWhenUnneeded=yes
