# systemd/clawpi-setup.target
[Unit]
Description=ClawPi setup target
AllowIsolate=yes
StopWhenUnneeded=yes
Wants=clawpi-setupd.service
Wants=clawpi-portald.service
