# systemd/clawpi-recovery.target
[Unit]
Description=ClawPi recovery target
AllowIsolate=yes
StopWhenUnneeded=yes
