# systemd/clawpi-portald.service
[Unit]
Description=ClawPi setup-network portal service
ConditionPathExists=/usr/local/libexec/clawpi/clawpi-portald
PartOf=clawpi-setup.target
After=local-fs.target clawpi-setupd.service

[Service]
Type=simple
ExecStart=/usr/local/libexec/clawpi/clawpi-portald
Restart=on-failure
RestartSec=2
