# systemd/clawpi-webd.service
[Unit]
Description=ClawPi local web landing service
ConditionPathExists=/usr/local/libexec/clawpi/clawpi-webd
After=local-fs.target clawpi-agentd.service avahi-daemon.service
Wants=clawpi-agentd.service
Wants=avahi-daemon.service

[Service]
Type=simple
ExecStart=/usr/local/libexec/clawpi/clawpi-webd
Restart=on-failure
RestartSec=2
